# lesson9/crontab
0 * * * * /opt/backup.sh
10 * * * * /opt/rsync_backup.sh